# srcs/ftps/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/supervisord.log
pidfile=/supervisord.pid

[program:vsftpd]
command=vsftpd /etc/vsftpd/vsftpd.conf
process_name=vsftpd
autostart=true
autorestart=false

[eventlistener:processes]
command=bash -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL
